# deploy/hermes-swarm.service
# systemd unit for a bare-metal Hermes Swarm deployment.
#
# Install:
#   1. Create a dedicated user + data dir:
#        sudo useradd --system --create-home --home-dir /var/lib/hermes-swarm hermes
#   2. Install the package into a venv owned by that user (or system-wide):
#        sudo -u hermes python3 -m venv /var/lib/hermes-swarm/venv
#        sudo -u hermes /var/lib/hermes-swarm/venv/bin/pip install /path/to/hermes-swarm
#        sudo -u hermes /var/lib/hermes-swarm/venv/bin/playwright install chromium
#   3. Put your config in /etc/hermes-swarm.env (see EnvironmentFile below).
#   4. sudo cp deploy/hermes-swarm.service /etc/systemd/system/
#      sudo systemctl daemon-reload && sudo systemctl enable --now hermes-swarm
#
# Verify the unit before enabling:  systemd-analyze verify deploy/hermes-swarm.service

[Unit]
Description=Hermes Swarm — multi-agent server + dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hermes
Group=hermes
# All runtime knobs (SWARM_API_KEY, SWARM_LLM_BASE_URL, SWARM_LLM_API_KEY,
# SWARM_DATA_DIR, SWARM_LOG_FILE, SWARM_HOST/PORT, …) live here. chmod 600 it.
EnvironmentFile=/etc/hermes-swarm.env
WorkingDirectory=/var/lib/hermes-swarm
ExecStart=/var/lib/hermes-swarm/venv/bin/hermes-swarm up
Restart=on-failure
RestartSec=5
# Graceful shutdown: the app's lifespan finalizer stops daemons + flushes
# browser cookies, so give it room before SIGKILL.
TimeoutStopSec=120
KillSignal=SIGTERM

# Hardening (agents can run terminal commands as this user — contain them).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hermes-swarm
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
